# infra/ansible/files/gai_model.service
[Unit]
Description=GAI Model Service
After=network.target

[Service]
Type=simple
ExecStart=/model/src/app.py
Restart=on-failure